# schoolDaemon.service
[Unit]
Description=School Daemon Service
After=network.target

[Service]
WorkingDirectory=/opt/schoolDaemon
Type=simple
Restart=always
User=root
Environment=XAUTHORITY=/home/student/.Xauthority
Environment=DISPLAY=:0
ExecStart=/opt/schoolDaemon/main.py

[Install]
WantedBy=multi-user.target
