# schoolDaemon.service
[Unit]
Description=School Daemon Service
After=network.target

[Service]
WorkingDirectory=/opt/schoolDaemon
Type=simple
Restart=always
User=root
ExecStart=/opt/schoolDaemon/main.py
ExecStartPre=/opt/schUpdater/updater.py

[Install]
WantedBy=multi-user.target
